# .docker/supervisord.conf
[supervisord]
nodaemon=true

[program:web-server]
directory=/app
command=php -S 0.0.0.0:80 -t ./public
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
